# templates/minecraft_server.socket.j2
[Socket]
ListenFIFO=%t/minecraft_server.pipe
Service=minecraft_server.service