# templates/minecraft_server.socket.j2
[Socket]
ListenFIFO=%t/minecraft_server.stdin
DeferAcceptSec=5
Service=minecraft_server.service