# templates/minecraft_server.socket.j2
[Socket]
ListenFIFO=%t/minecraft_server.stdin
Service=minecraft_server.service